# crontab/run-jupyter
#!/bin/bash
jupyter notebook
jupyter nbextension enable --py widgetsnbextension
